# lab3/wordpress/supervisord-mysqld.conf
[program:mysqld]
command=/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true

[program:cron]
command=/usr/sbin/cron
numprocs=1
autostart=true
autorestart=true
